# systemd/system/vcsFPService.service.d/local.conf
[Service]
PrivateTmp=true
JoinsNamespaceOf=fprintd.service
